# templates/paperless-flower.service.j2
[Unit]
Description=Paperless Celery Flower
Requires=redis.service

[Service]
User=tbd
Group=tbd
WorkingDirectory=tbd
ExecStart=tbd

[Install]
WantedBy=multi-user.target